# rootskel/src/etc/inittab
# /etc/inittab
# busybox init configuration for debian-installer

# main rc script
null::sysinit:/etc/init.d/rcS

# main setup program
::respawn:/sbin/debian-installer

# convenience shell
vc/2::askfirst:/bin/sh

# logging
::sysinit:/bin/sh -c "[ -e /dev/vc/4 ] && /bin/dmesg > /dev/vc/4"
vc/3::respawn:/usr/bin/tail -f /var/log/messages
vc/4::respawn:/usr/bin/tail -f /proc/kmsg

# Stuff to do before rebooting
::ctrlaltdel:/bin/umount -a -r > /dev/null 2>&1
::ctrlaltdel:/sbin/swapoff -a > /dev/null 2>&1

