# packages/rootskel/src/etc/inittab
# /etc/inittab
# busybox init configuration for debian-installer

# main rc script
::sysinit:/sbin/debian-installer-startup

# main setup program
::respawn:/sbin/debian-installer

# convenience shell
vc/2::askfirst:/bin/sh

# logging
vc/3::respawn:/usr/bin/tail -f /var/log/messages
vc/4::respawn:/usr/bin/tail -f /var/log/syslog

# Stuff to do before rebooting
::ctrlaltdel:/sbin/shutdown > /dev/null 2>&1

